# cms/explore-dubuque-cms.service
[Unit]
Description=Explore Dubuque CMS (Payload 3 / Next.js)
After=network.target postgresql.service

[Service]
Type=simple
User=explore
WorkingDirectory=/explore/cms
ExecStart=/usr/bin/node node_modules/.bin/next start --port 3001
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=explore-dubuque-cms
EnvironmentFile=/explore/cms/.env

[Install]
WantedBy=multi-user.target
